# rpm/src/main/rpm/etc/profile.d/archiva.sh
export RUN_AS_USER=archiva
export ARCHIVA_BASE=/var/lib/archiva
